# systemd/netcfg-netns.service
[Unit]
Description=netcfg network namespace configuration
Documentation=https://github.com/netcfg/netcfg
After=netcfg.service network.target
Wants=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/netcfg apply --netns-only
ExecStop=/usr/bin/netcfg destroy

[Install]
WantedBy=multi-user.target
